# configs/bootstrap/consul-bootstrap.conf
description "Consul bootstrap process"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]
respawn 

console log
exec consul agent -ui -data-dir /var/consul/ -config-dir /etc/consul.d/bootstrap -encrypt ${encrypt} -bind ${bindip}
